# computer/etc/linux-systemd/user/parachute.service
[Unit]
Description=Parachute Base Server
Documentation=https://github.com/OpenParachutePBC/parachute
After=network.target

[Service]
Type=simple
Environment=VAULT_PATH=%h/Parachute
Environment=PORT=3333
Environment=HOST=0.0.0.0
ExecStart=/usr/local/bin/parachute-server
Restart=on-failure
RestartSec=10
SuccessExitStatus=143

[Install]
WantedBy=default.target
