# computer/etc/macos-launchd/com.parachute.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.parachute.server</string>

    <key>Comment</key>
    <string>Parachute Base Server - AI orchestration for your knowledge vault</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/parachute-server</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>VAULT_PATH</key>
        <string>~/Parachute</string>
        <key>PORT</key>
        <string>3333</string>
        <key>HOST</key>
        <string>0.0.0.0</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/parachute.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/parachute.log</string>

    <key>WorkingDirectory</key>
    <string>/usr/local/var/parachute</string>
</dict>
</plist>
